# deploy/embeddings_service/zeta-embeddings.service
[Unit]
Description=Zeta Embeddings Service (Docker Compose)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/zeta/CHATBOT2.0/deploy/embeddings_service
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
